# clean-workspaces.service
[Unit]
Description="Clean up expired workspaces"

[Service]
type=oneshot
ExecStart=/usr/local/bin/workspaces clean
