# workspace-clean.timer
[Unit]
Description=Clean up expired workspaces daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
